# etc/sv/amneziavpn/run
#!/bin/sh
exec 2>&1
exec /opt/AmneziaVPN/service/./AmneziaVPN-service.sh
