# etc/sv/amneziavpn/run
#!/bin/sh
exec 2>&1
exec /home/snake/.local/state/nix/profile/bin/AmneziaVPN-service
