# etc/profile.d/sdkman.sh
#!/bin/bash

export SDKMAN_DIR="/home/snake/.local/devjava/sdkman"


#[[ -s "$SDKMAN_DIR/bin/sdkman-init.sh" ]] && source "$SDKMAN_DIR/bin/sdkman-init.sh"
