# etc/sv/runsvdir-snake/run
#!/bin/sh

export USER="snake"
export HOME="/home/snake"

groups="$(id -Gn "$USER" | tr ' ' ':')"
svdir="$HOME/.local/services"

exec chpst -u "$USER:$groups" runsvdir "$svdir"

